# bots/meic/meic.service
# MEIC (Multiple Entry Iron Condors) Trading Bot Systemd Service
#
# Installation:
#   sudo cp /opt/calypso/bots/meic/meic.service /etc/systemd/system/meic.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable meic.service
#   sudo systemctl start meic.service
#
# Management:
#   sudo systemctl status meic
#   sudo systemctl stop meic
#   sudo systemctl restart meic
#   sudo journalctl -u meic -f
#
# IMPORTANT: Never use `kill` or `pkill` - always use systemctl to stop the service!
# The service has Restart=always, so killing it will cause auto-restart in 30 seconds.

[Unit]
Description=MEIC 0DTE Trading Bot (Tammy Chambless Strategy)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=calypso
Group=calypso
WorkingDirectory=/opt/calypso
Environment="PATH=/opt/calypso/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/calypso"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/opt/calypso/credentials/service-account.json"

# Run the MEIC bot
# Mode is controlled by config.json "dry_run" field:
#   "dry_run": true  = simulation (no real orders)
#   "dry_run": false = live trading
# CLI flag --dry-run also works and takes priority over config
ExecStart=/opt/calypso/.venv/bin/python -m bots.meic.main

# Graceful shutdown - send SIGTERM, wait 60s for cleanup
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Restart policy
# Bot will automatically restart on crashes
# Max 5 restarts in 10 minutes, then fail
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=5

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meic

[Install]
WantedBy=multi-user.target
