# bots/hydra/hydra.service
# HYDRA 0DTE Trading Bot Systemd Service
#
# Multi-Entry Iron Condors (SPX 0DTE) with credit gates, progressive OTM
# tightening, and smart exits. Based on Tammy Chambless's MEIC strategy.
#
# Installation:
#   sudo cp /opt/calypso/bots/hydra/hydra.service /etc/systemd/system/hydra.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable hydra.service
#   sudo systemctl start hydra.service
#
# Management:
#   sudo systemctl status hydra
#   sudo systemctl stop hydra
#   sudo systemctl restart hydra
#   sudo journalctl -u hydra -f
#
# IMPORTANT: Never use `kill` or `pkill` - always use systemctl to stop the service!
# The service has Restart=always, so killing it will cause auto-restart in 30 seconds.

[Unit]
Description=HYDRA 0DTE Trading Bot (Multi-Entry Iron Condors)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=calypso
Group=calypso
WorkingDirectory=/opt/calypso

# Run the HYDRA bot
# Mode is controlled by config.json "dry_run" field:
#   "dry_run": true  = simulation (no real orders)
#   "dry_run": false = live trading
# CLI flag --dry-run also works and takes priority over config
ExecStart=/opt/calypso/.venv/bin/python -m bots.hydra.main

# Graceful shutdown - send SIGTERM, wait 60s for cleanup
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Restart policy
# Bot will automatically restart on crashes
# Max 5 restarts in 10 minutes, then fail
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hydra

[Install]
WantedBy=multi-user.target
