# deploy/homer.service
[Unit]
Description=Calypso HOMER Trading Journal Writer (oneshot)
Documentation=https://github.com/calypso-trading/calypso

[Service]
Type=oneshot
User=calypso
Group=calypso
WorkingDirectory=/opt/calypso

Environment="GCP_PROJECT=calypso-trading-bot"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/calypso"

ExecStart=/opt/calypso/.venv/bin/python -m services.homer.main

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=calypso-homer

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# HOMER does `git pull --rebase` which may touch any file in /opt/calypso that
# has changed on origin (e.g. frontend src files, strategy modules). Limiting
# ReadWritePaths to specific subdirectories caused silent `git rebase` failures
# with "Read-only file system" (observed 2026-04-21 when the VM was behind
# origin by frontend commits). /opt/calypso is HOMER's workspace — the whole
# tree must be writable. ProtectSystem=strict still shields the rest of /.
ReadWritePaths=/opt/calypso

# Claude API (120s per narrative × 3 calls per day × N days) + Sheets reads (3×15s) + git (60s) + overhead
TimeoutStartSec=600
