# deploy/hydra_variant_b.service
[Unit]
Description=HYDRA Variant B (110pt spread, dry-run, parallel head-to-head)
After=network.target token_keeper.service
Wants=token_keeper.service

[Service]
Type=simple
User=calypso
Group=calypso
WorkingDirectory=/opt/calypso
Environment="HYDRA_VARIANT_ID=b"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/calypso/.venv/bin/python -m bots.hydra.main --config bots/hydra/config/config_variant_b.json
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hydra_variant_b
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
